# .config/rc/init.d/downtime
#!/sbin/openrc-run


start() {
    ebegin "Checking website availability"
    if ! curl --fail -s -o /dev/null http://192.168.100.245/login.html; then
        notify-send "Website Down" "Could not access http://192.168.100.245/login.html"
    fi
    eend $?
}
