# .config/rc/init.d/prostat
#!/sbin/openrc-run

description="mount YEET"

depend() {
    keyword -prefix
    need dbus
}

start() {
    ebegin "mount Y E E T"
mkdir -p /run/media/shigure/YEET
udisksctl mount -b /dev/disk/by-id/usb-JMicron_Generic_0123456789ABCDEF-0:0-part1 -t exfat /run/media/shigure/YEET
chown shigure:shigure /run/media/shigure/YEET
    eend $?
}
