# .config/rc/init.d/yandex
#!/sbin/openrc-run

description="Yandex.Disk client daemon"
command="/usr/bin/yandex-disk"
if [ -n "${RC_USER}" ]; then
	pidfile="${XDG_RUNTIME_DIR:-/run/user/$(id -u)}/${SVCNAME}.pid"
else
	pidfile="/run/${SVCNAME}.pid"
fi

# This script is based on the systemd service file provided.
# The Restart=always and RestartSec=30 from the systemd unit are not
# directly translated. OpenRC can handle this with a supervisor.
# If you are using a supervisor like supervise-daemon, you can add
# 'supervisor supervise-daemon' to this script and set the appropriate options.

# The configuration file path is taken from the systemd unit.
# You can override this by creating /etc/conf.d/yandex and setting
# YANDEX_DISK_CONFIG_FILE="/path/to/your/config.cfg"
: ${YANDEX_DISK_CONFIG_FILE:="${HOME}/.config/yandex-disk/config.cfg"}


start_pre() {
	if [ -n "${RC_USER}" ]; then
		checkpath -d -o "${RC_USER}" "$(dirname "${pidfile}")"
	fi
}

start() {
    ebegin "Starting Yandex.Disk"
    start-stop-daemon --start --quiet --background \
        --make-pidfile --pidfile "${pidfile}" \
        --exec "${command}" -- start -c "${YANDEX_DISK_CONFIG_FILE}"
    eend $?
}

stop() {
    ebegin "Stopping Yandex.Disk"
    start-stop-daemon --stop --quiet --pidfile "${pidfile}"
    eend $?
}

reload() {
    ebegin "Reloading Yandex.Disk"
    /usr/bin/killall -qw yandex-disk
    eend $?
}

