# .config/rc/init.d/YEET
#!/sbin/openrc-run

description="mount YEET"

command="/usr/bin/udisksctl"
command_args="mount -b /dev/disk/by-id/usb-JMicron_Generic_0123456789ABCDEF-0:0-part1 -t exfat /run/media/shigure/YEET"
command_user="shigure"

supervisor=supervise-daemon
#supervise_daemon_args="--oneshot"

depend() {
    need dbus
    keyword -prefix
}
